# supervisord.conf
[supervisord]
nodaemon=true

[program:shadowvpn_server]
command=/usr/local/bin/shadowvpn -c /etc/shadowvpn/server.conf
startsecs=10
startretries=3
autostart=true
autorestart=true
